# overlay/etc/service/pocoirc/run
#!/bin/sh
cd /usr/src/app
exec pocoirc --verbose --config config.yaml
