# overlay/etc/service/servizietto/run
#!/bin/sh
# put your servizietto here
# remember: it can't daemonize

# exec echo "Servizietto's placeholder"
