# overlay/etc/service/bot/run
#!/bin/sh
# put your servizietto here
# remember: it can't daemonize
# more details http://smarden.org/runit/runsv.8.html

cd /usr/src/bot
exec pocoirc --verbose --config etc/config.yaml
